# lima/veronica.service
[Unit]
Description=Veronica eBPF Intelligence Layer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/veronicad
Environment=AGENTFIELD_URL=http://host.lima.internal:8090
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
