# lima/veronica.service
[Unit]
Description=Veronica eBPF Intelligence Layer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/veronicad
Environment=VERONICA_STORE_DIR=/var/veronica/nats
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
